# ai-trading-scheduler.service
[Unit]
Description=AI Trading Bot Scheduler & Metrics
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ai-trading-bot
ExecStart=/opt/ai-trading-bot/start.sh
Restart=always
RestartSec=5
StandardOutput=append:/var/log/ai-trading-scheduler.log
StandardError=append:/var/log/ai-trading-scheduler.log

[Install]
WantedBy=multi-user.target
