# packaging/systemd/ai-trading-research-weekend-sunday.timer
[Unit]
Description=Run AI Trading Weekend Sunday Validation Automation

[Timer]
OnCalendar=Sun 10:00 America/New_York
RandomizedDelaySec=10m
AccuracySec=1m
Persistent=true
Unit=ai-trading-research-weekend-sunday.service

[Install]
WantedBy=timers.target
