# packaging/systemd/ai-trading-healthcheck.service
[Unit]
Description=AI Trading Runtime Health Check
After=network-online.target ai-trading.service
Wants=network-online.target

[Service]
Type=oneshot
User=aiuser
Group=aiuser
WorkingDirectory=/home/aiuser/ai-trading-bot
Environment=PATH=/home/aiuser/ai-trading-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=AI_TRADING_ENV_SRC=/home/aiuser/ai-trading-bot/.env
Environment=AI_TRADING_RUNTIME_ENV_DST=/run/ai-trading-bot/ai-trading-runtime.env
Environment=AI_TRADING_RUNTIME_ENV_PATH=/run/ai-trading-bot/ai-trading-runtime.env
Environment=HEALTHCHECK_URL=http://127.0.0.1:9001/healthz
Environment=HEALTHCHECK_SYSTEMD_UNIT=ai-trading.service
EnvironmentFile=-/run/ai-trading-bot/ai-trading-runtime.env
Environment=AI_TRADING_DOTENV_RUNTIME_OVERRIDE=1
ExecStartPre=/home/aiuser/ai-trading-bot/scripts/sync_env_runtime.sh
ExecStart=/bin/bash -lc 'bash /home/aiuser/ai-trading-bot/scripts/runtime_phase1_health_check_runner.sh'
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
RuntimeDirectory=ai-trading-bot
ReadWritePaths=/run/ai-trading-bot

[Install]
WantedBy=multi-user.target
