# packaging/systemd/ai-trading-runtime-prune.service
[Unit]
Description=Prune AI Trading Runtime JSONL Artifacts
After=ai-trading.service
OnFailure=ai-trading-systemd-failure-alert@%p.service

[Service]
Type=oneshot
User=aiuser
Group=aiuser
WorkingDirectory=/home/aiuser/ai-trading-bot
ExecStart=/bin/bash -lc 'bash /home/aiuser/ai-trading-bot/scripts/prune_runtime_jsonl.sh'
