# packaging/systemd/ai-trading-connectors.service
[Unit]
Description=Dispatch AI Trading Incident Connectors (Slack + OpenClaw)
After=network-online.target ai-trading.service
Wants=network-online.target

[Service]
Type=oneshot
User=aiuser
Group=aiuser
WorkingDirectory=/home/aiuser/ai-trading-bot
Environment=PATH=/home/aiuser/ai-trading-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=AI_TRADING_ENV_SRC=/home/aiuser/ai-trading-bot/.env
Environment=AI_TRADING_RUNTIME_ENV_DST=/run/ai-trading-bot/ai-trading-runtime.env
Environment=AI_TRADING_RUNTIME_ENV_PATH=/run/ai-trading-bot/ai-trading-runtime.env
EnvironmentFile=-/run/ai-trading-bot/ai-trading-runtime.env
Environment=AI_TRADING_DOTENV_RUNTIME_OVERRIDE=1
Environment=API_PORT=9001
Environment=HEALTHCHECK_PORT=9001
ExecStartPre=/home/aiuser/ai-trading-bot/scripts/sync_env_runtime.sh
ExecStart=/home/aiuser/ai-trading-bot/venv/bin/python /home/aiuser/ai-trading-bot/scripts/connector_incident_dispatch.py
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
RuntimeDirectory=ai-trading-bot
ReadWritePaths=/var/lib/ai-trading-bot/runtime /run/ai-trading-bot

[Install]
WantedBy=multi-user.target
