# packaging/systemd/ai-trading-systemd-failure-alert@.service
[Unit]
Description=Notify on failed ai-trading systemd unit (%I)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=aiuser
Group=aiuser
WorkingDirectory=/home/aiuser/ai-trading-bot
EnvironmentFile=-/home/aiuser/ai-trading-bot/.env
EnvironmentFile=-/home/aiuser/ai-trading-bot/.env.runtime
ExecStart=/bin/bash -lc 'bash /home/aiuser/ai-trading-bot/scripts/systemd_failure_alert.sh "%I"'

