# packaging/systemd/ai-trading-runtime-report.service
[Unit]
Description=Refresh AI Trading Runtime Performance Reports
After=network-online.target ai-trading.service
Wants=network-online.target
OnFailure=ai-trading-systemd-failure-alert@%p.service

[Service]
Type=oneshot
User=aiuser
Group=aiuser
WorkingDirectory=/home/aiuser/ai-trading-bot
Environment=PATH=/home/aiuser/ai-trading-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=AI_TRADING_DATA_DIR=/var/lib/ai-trading-bot
Environment=AI_TRADING_CACHE_DIR=/var/cache/ai-trading-bot
Environment=AI_TRADING_LOG_DIR=/var/log/ai-trading-bot
EnvironmentFile=-/home/aiuser/ai-trading-bot/.env.runtime
EnvironmentFile=-/home/aiuser/ai-trading-bot/.env
ExecStart=/bin/bash -lc 'bash /home/aiuser/ai-trading-bot/scripts/refresh_runtime_reports.sh'
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/ai-trading-bot /var/cache/ai-trading-bot /var/log/ai-trading-bot

[Install]
WantedBy=multi-user.target
