# packaging/systemd/ai-trading.service
[Unit]
Description=AI Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aiuser
Group=aiuser
WorkingDirectory=/home/aiuser/ai-trading-bot
Environment=PATH=/home/aiuser/ai-trading-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=ALPACA_TRADING_BASE_URL=https://paper-api.alpaca.markets
Environment=AI_TRADING_ENV_SRC=/etc/ai-trading-bot/ai-trading.env
Environment=AI_TRADING_RUNTIME_ENV_DST=/run/ai-trading-bot/ai-trading-runtime.env
Environment=AI_TRADING_RUNTIME_ENV_PATH=/run/ai-trading-bot/ai-trading-runtime.env
EnvironmentFile=-/etc/ai-trading-bot/ai-trading.env
EnvironmentFile=-/run/ai-trading-bot/ai-trading-runtime.env
Environment=AI_TRADING_MODEL_MODULE=ai_trading.model_loader
Environment=AI_TRADING_MODEL_PATH=/var/lib/ai-trading-bot/models/trained_model.pkl
Environment=AI_TRADING_DATA_DIR=/var/lib/ai-trading-bot
Environment=AI_TRADING_CACHE_DIR=/var/cache/ai-trading-bot
Environment=AI_TRADING_LOG_DIR=/var/log/ai-trading-bot
Environment=AI_TRADING_MODELS_DIR=/var/lib/ai-trading-bot/models
Environment=AI_TRADING_OUTPUT_DIR=/var/lib/ai-trading-bot/output
Environment=API_PORT=9001
# Match HEALTHCHECK_PORT to API_PORT so the shared Flask process serves /healthz on :9001
Environment=HEALTHCHECK_PORT=9001
Environment=TRADE_LOG_PATH=/var/log/ai-trading-bot/trades.jsonl
StateDirectory=ai-trading-bot
CacheDirectory=ai-trading-bot
LogsDirectory=ai-trading-bot
RuntimeDirectory=ai-trading-bot
ExecStartPre=/usr/bin/install -d -m 700 \
    /var/lib/ai-trading-bot/models \
    /var/lib/ai-trading-bot/output \
    /var/lib/ai-trading-bot/runtime
ExecStartPre=/usr/bin/install -d -m 775 \
    /var/lib/ai-trading-bot/runtime/research_reports
ExecStartPre=/usr/bin/find /var/lib/ai-trading-bot/runtime/research_reports -type f \
    -name after_hours_training*.json -exec /usr/bin/chmod 664 {} +
ExecStartPre=/home/aiuser/ai-trading-bot/scripts/sync_env_runtime.sh
ExecStartPre=/home/aiuser/ai-trading-bot/venv/bin/alembic upgrade head
ExecStart=/home/aiuser/ai-trading-bot/venv/bin/python -m ai_trading
# Allow long-lived process; systemd will supervise restarts
RuntimeMaxSec=infinity
Restart=always
RestartSec=10
RestartPreventExitStatus=98
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/ai-trading-bot /var/cache/ai-trading-bot /var/log/ai-trading-bot /run/ai-trading-bot

[Install]
WantedBy=multi-user.target
