# deploy/ai-trading.service
[Unit]
Description=AI Trading Bot Scheduler
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=/home/aiuser/ai-trading-bot
ExecStart=/home/aiuser/ai-trading-bot/start.sh
Environment=BOT_LOG_FILE=/home/aiuser/ai-trading-bot/logs/scheduler.log
Environment=SCHEDULER_SLEEP_SECONDS=30
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
